# cups-auto-configure.service
[Unit]
Description=Auto-configure CUPS printer from .env
After=cups.service
Requires=cups.service

[Service]
Type=oneshot
ExecStart=%PROJECT_DIR%/installers/cups-auto-configure.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target

