# rpi-waybill-printer-backend.service
[Unit]
Description=RPI Waybill Printer Backend (Flask API)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%PROJECT_DIR%
Environment="PATH=%PROJECT_DIR%/venv/bin"
Environment="FLASK_APP=run:app"
Environment="FLASK_ENV=production"
ExecStart=%PROJECT_DIR%/venv/bin/python run.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

