# rpi-waybill-printer-frontend.service
[Unit]
Description=RPI Waybill Printer Frontend (React + Vite Dev Server)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%PROJECT_DIR%/frontend
Environment="PATH=%PROJECT_DIR%/node_modules/.bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/bin/bash -c "cd %PROJECT_DIR%/frontend && npm run dev -- --host 0.0.0.0"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
