# systemd/rpi-waybill-printer.service
[Unit]
Description=RPI Waybill Printer - Flask Backend + React Frontend
After=network.target
Wants=rpi-waybill-printer-frontend.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rpi-waybill-printer
Environment="PATH=/home/pi/rpi-waybill-printer/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="USE_HTTPS=true"
ExecStart=/home/pi/rpi-waybill-printer/venv/bin/python3 /home/pi/rpi-waybill-printer/run.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

